# blog.conf
[program:blog]
command=python3 /usr/local/bin/gunicorn wsgi -c gunicorn.config.py
directory=/root/deploy/blog
autostart=true
autorestart=true
